# video-zfajax-test.conf
[program:video-zfajax-test]
process_name=%(program_name)s
#需要修改路径
command=php <你项目的路径>/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
#需要修改路径
stdout_logfile=<你项目的路径>/storage/logs/worker.log
stopwaitsecs=3600
